# hiltop_board/rootfs_overlay/etc/udev/rules.d/10-rtc.rules
SUBSYSTEM=="rtc", ACTION=="add", RUN+="/sbin/hwclock -s -f /dev/%k" 
