# hiltop_board/rootfs_overlay/etc/network/if-up.d/start_hwclock_timer
#! /bin/sh
at now + 1 hour -f /opt/devtank/set_hwclock.sh
